# dell2-admin-maas/dell_var_lib_tftpboot_centos7_x64_everything_1804_scripts/after-install.service
[Unit]
Description=AfterInstallNVIDIA
Requires=network-online.target
After=network-online.target

[Service]
ExecStart=/bin/bash /opt/purkkaskripti.sh

[Install]
WantedBy=default.target
